# process-reaper.service
[Unit]
Description=Intelligent Process Reaper
Documentation=https://github.com/clawbotari/process-reaper
After=network.target

[Service]
Type=simple
User=root
Group=root
# Environment variables define the reaper's behavior
Environment=REAPER_PATTERN=.*myapp.*
Environment=REAPER_INTERVAL=60
Environment=REAPER_LOG_DIR=/var/log/process-reaper
Environment=REAPER_GRACE_PERIOD=10
Environment=REAPER_MIN_UPTIME=5
# UniVerse deep forensic integration (optional)
Environment=REAPER_UV_DIR=
# Audit mode: if false, only log & collect forensic data
Environment=REAPER_KILL=true
# Suppress heartbeat logs when no candidates are found
Environment=REAPER_HEARTBEAT_QUIET=false
ExecStart=/usr/local/bin/process-reaper
Restart=on-failure
RestartSec=10
StandardOutput=file:/var/log/process-reaper/process-reaper-service.log
StandardError=file:/var/log/process-reaper/process-reaper-service.log
SyslogIdentifier=process-reaper

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/log/process-reaper

[Install]
WantedBy=multi-user.target
